# templates/luma.socket.j2
[Unit]
Description={{luma_project}} socket

[Socket]
ListenStream=/run/luma_{{luma_project}}.sock

[Install]
WantedBy=sockets.target
